# es/modules/guide-basic-dns-whitehole-nfqueue/examples/etc/systemd/system/luids-xlist-autoremove.timer
[Unit]
Description=Truncate file autoremove.xlist every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
