# es/modules/guide-basic-dns-whitehole-nfqueue/examples/etc/systemd/system/luids-xlist-autoremove.service
[Unit]
Description=Truncate file autoremove.xlist

[Service]
Type=oneshot
User=root
ExecStart=truncate -s 0 /var/lib/luids/xlist/local/autoremove.xlist

